# drc.service.in
[Unit]
Description=DRC restore/save (user service)
After=default.target

# Template — generated to drc.service by install.sh from config.env.

[Service]
Type=oneshot
RemainAfterExit=yes

ExecStart=/bin/bash -c '\
  ARG="$(cat @REPO_DIR@/last_arg 2>/dev/null || echo off)"; \
  echo "[DRC] Restoring last_arg=$ARG"; \
  exec @REPO_DIR@/drc.sh "$ARG"'

[Install]
WantedBy=default.target
