# etc/rc.d/brutefir_drc.in
#!/bin/sh
#
# PROVIDE: brutefir_drc
# REQUIRE: musicpd
# KEYWORD: shutdown
#
# Template — generated to brutefir_drc by install.sh from config.env.
#
# Add to /etc/rc.conf:
# brutefir_drc_enable="YES"
#
# IMPORTANT: brutefir must run as the AUDIO USER, not root.  Both entry points
# (this boot/hotplug service AND an interactive `drc.sh`) have to own brutefir as
# the same user, otherwise a user-run `drc.sh` cannot kill a root-owned brutefir.
# Started as root at boot we drop privileges with su(1); started by the user
# already (e.g. `service brutefir_drc onestart`) we run in place.  virtual_oss
# still becomes root via drc.sh's internal sudo, so the loopback stays root
# infrastructure while brutefir stays user-owned.

. /etc/rc.subr

name="brutefir_drc"
rcvar="brutefir_drc_enable"
start_cmd="brutefir_drc_start"
stop_cmd="brutefir_drc_stop"
status_cmd="brutefir_drc_status"

load_rc_config "$name"

: ${brutefir_drc_enable:="NO"}
: ${brutefir_drc_user:="@AUDIO_USER@"}
: ${brutefir_drc_drcsh:="@REPO_DIR@/drc.sh"}
: ${brutefir_drc_statussh:="@REPO_DIR@/drc-status.sh"}

# Run a command as the audio user.  Use a LOGIN shell (su -l) so HOME and the
# full PATH (incl. /usr/local/bin, where brutefir/mpc/virtual_oss live) are set
# from the user's login class — a boot service's PATH would otherwise miss them.
_drc_as_user()
{
	if [ "$(id -u)" -eq 0 ]; then
		su -l "${brutefir_drc_user}" -c "$*"
	else
		sh -c "$*"
	fi
}

brutefir_drc_start()
{
	_drc_as_user "${brutefir_drc_drcsh} restore"
}

brutefir_drc_stop()
{
	_drc_as_user "${brutefir_drc_drcsh} off"
}

brutefir_drc_status()
{
	_drc_as_user "${brutefir_drc_statussh}"
}

run_rc_command "$1"
