# etc/rc.d/drc_usb_audio
#!/bin/sh
#
# PROVIDE: drc_usb_audio
# REQUIRE: brutefir_drc musicpd
# KEYWORD: shutdown
#
# Add to /etc/rc.conf:
# drc_usb_audio_enable="YES"
#
# This is the single entry point that brings DRC up when the DAC is present.
# It is driven from two places that funnel through the same idempotent start:
#   - boot: rc starts it once; it PROBES for the DAC (level check) and acts only
#           if the DAC is already on.  If the DAC is off, it does nothing and
#           defers to devd.
#   - hotplug: the devd ATTACH rule runs `onestart` when a USB audio interface
#           appears (DAC switched on later).
# The devd DETACH rule runs `onestop` so the next attach re-triggers cleanly.
# brutefir_drc is NOT enabled as a standalone boot service anymore — it is just
# the worker invoked from here (and it runs `drc.sh restore`).
#

. /etc/rc.subr

name="drc_usb_audio"
rcvar="drc_usb_audio_enable"
start_cmd="drc_usb_audio_start"
stop_cmd="drc_usb_audio_stop"
status_cmd="drc_usb_audio_status"

load_rc_config "$name"

: ${drc_usb_audio_enable:="NO"}
: ${drc_usb_audio_start_delay:="1"}
: ${drc_usb_audio_state_file:="/var/run/drc_usb_audio.active"}
: ${drc_usb_audio_brutefir_service:="brutefir_drc"}
# OSS node brutefir opens; its presence is our proxy for "the DAC is here".
: ${drc_usb_audio_dsp_dev:="/dev/dsp0"}

drc_usb_audio_start()
{
	if [ -f "${drc_usb_audio_state_file}" ]; then
		echo "${name} already active"
		return 0
	fi

	# Settle first: on hotplug devd fires when the USB interface attaches,
	# a moment before the OSS node is created — so wait, THEN probe.  At boot
	# the node already exists and the wait is harmless.
	sleep "${drc_usb_audio_start_delay}"

	# Only act when the DAC is actually present.  If it is not, do nothing —
	# the devd ATTACH rule will start us when the DAC appears.  This is what
	# lets us drop the unconditional brutefir_drc boot service.
	if [ ! -e "${drc_usb_audio_dsp_dev}" ]; then
		echo "${name}: ${drc_usb_audio_dsp_dev} not present; deferring to devd"
		return 0
	fi

	service "${drc_usb_audio_brutefir_service}" onestart || return 1
	: > "${drc_usb_audio_state_file}"
}

drc_usb_audio_stop()
{
	service "${drc_usb_audio_brutefir_service}" onestop
	rm -f "${drc_usb_audio_state_file}"
}

drc_usb_audio_status()
{
	if [ -f "${drc_usb_audio_state_file}" ]; then
		echo "${name} active"
		service "${drc_usb_audio_brutefir_service}" onestatus
		return 0
	fi

	echo "${name} inactive"
	return 1
}

run_rc_command "$1"
