# etc/rc.d/musicpd.in
#!/bin/sh

# PROVIDE: musicpd
# REQUIRE: mixer LOGIN avahi_daemon oss
# KEYWORD: shutdown

# Template — generated to musicpd by install.sh from config.env.
#
# Add the following line to /etc/rc.conf to enable mpd:
#
# musicpd_enable=YES
#
# Run-from-repo model: this script reads the config DIRECTLY from the git
# checkout instead of a copy under /usr/local/etc, so there is nothing to
# re-copy on update (git pull is enough).  mpd is started by root and drops to
# the unprivileged "user"/"group" named inside the config file.
#
# Override musicpd_config in /etc/rc.conf (or /etc/rc.conf.d/musicpd) if the
# checkout lives elsewhere.

. /etc/rc.subr

name=musicpd
rcvar=musicpd_enable
desc="Music Player Daemon"

load_rc_config $name

: ${musicpd_enable:=NO}
: ${musicpd_config:="@REPO_DIR@/mpd/musicpd.conf"}

command=/usr/local/bin/musicpd
command_args="${musicpd_config}"
required_files="${musicpd_config}"
start_precmd=${name}_getpidfile
stop_precmd=${name}_getpidfile

musicpd_getpidfile()
{
	if get_pidfile_from_conf pid_file "${musicpd_config}" ; then
		pidfile="$_pidfile_from_conf"
	else
		pidfile="@AUDIO_HOME@/.cache/mpd/pid"
	fi
}

run_rc_command $1
