# etc/systemd/system/mpd.service.d/open-media-drc.conf
# Drop-in generated by install.sh from config.env.
# Overrides /usr/lib/systemd/system/mpd.service.d/00-arch.conf (User=mpd)
# and adds the repo config path — /etc/ drop-ins take precedence over
# /usr/lib/ drop-ins regardless of filename order.
# Deploy to /etc/systemd/system/mpd.service.d/omdrc.conf

[Service]
User=giacomo
ExecStart=
ExecStart=/usr/bin/mpd --systemd /home/giacomo/open-media-drc/mpd/mpd.conf
