# etc/systemd/user/upmpdcli.service.in
[Unit]
Description=upmpdcli UPnP/OpenHome MPD renderer (user service)
After=network-online.target sound.target
Wants=network-online.target

# Template — generated to upmpdcli.service by install.sh from config.env.
# systemd --user unit: runs as the logged-in audio user (no User= line), reads
# its config straight from the checkout.  Install to ~/.config/systemd/user/ and
# enable lingering so it survives logout:  loginctl enable-linger @AUDIO_USER@

[Service]
Type=simple
ExecStart=@PREFIX@/bin/upmpdcli -c @REPO_DIR@/upmpdcli/upmpdcli.conf
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
