# etc/systemd/system/brutefir-drc.service.in
[Unit]
Description=BruteFIR DRC convolution engine
After=sound.target

# Template — generated to brutefir-drc.service by install.sh from config.env.
# Runs as the AUDIO USER (not root): brutefir must be owned by that user so an
# interactive drc.sh (run as the same user) can stop/restart it.  HOME points at
# the audio user so brutefir reads its XDG defaults from @AUDIO_HOME@/.config/
# BruteFIR.  LimitRTPRIO/LimitMEMLOCK let brutefir set realtime scheduling and
# lock memory while unprivileged (otherwise it warns and runs non-RT).

[Service]
Type=oneshot
RemainAfterExit=yes
User=@AUDIO_USER@
Environment=HOME=@AUDIO_HOME@
Environment=PATH=/usr/local/bin:/usr/bin:/bin
LimitRTPRIO=99
LimitMEMLOCK=infinity
ExecStart=@REPO_DIR@/drc.sh restore
ExecStop=@REPO_DIR@/drc.sh off
