# etc/systemd/system/drc-usb-audio.service.in
[Unit]
Description=DRC audio routing when USB DAC is connected
After=sound.target

# Template — generated to drc-usb-audio.service by install.sh from config.env.

[Service]
Type=oneshot
RemainAfterExit=yes
User=@AUDIO_USER@
Environment=HOME=@AUDIO_HOME@
Environment=PATH=/usr/local/bin:/usr/bin:/bin
LimitRTPRIO=99
LimitMEMLOCK=infinity
ExecStartPre=/bin/sleep 1
ExecStart=@REPO_DIR@/drc.sh restore
ExecStop=@REPO_DIR@/drc.sh off
