# etc/systemd/system/mpd.service.in
[Unit]
Description=Music Player Daemon
Documentation=man:mpd(1) man:mpd.conf(5)
After=network.target sound.target

# Template — generated to mpd.service by install.sh from config.env.
# Run-from-repo model: MPD runs as @AUDIO_USER@ and reads its config directly
# from this checkout, mirroring the FreeBSD musicpd rc.d model.
# Deploy to /etc/systemd/system/mpd.service (shadows the package unit), then:
#   sudo systemctl disable --now mpd.socket   # disable package socket activation
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now mpd.service
# State dirs must exist first:
#   mkdir -p @AUDIO_HOME@/.local/share/mpd @AUDIO_HOME@/.cache/mpd

[Service]
Type=notify
User=@AUDIO_USER@
ExecStart=/usr/bin/mpd --systemd @REPO_DIR@/mpd/mpd.conf
LimitRTPRIO=40
LimitRTTIME=infinity
LimitMEMLOCK=64M
ProtectSystem=yes
NoNewPrivileges=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes

[Install]
WantedBy=multi-user.target
